# demonstration/manufacture.service
[Unit]
Description=Angle manufacture

[Service]
ExecStart=/etc/systemd/system/start_manufacture.sh

[Install]
WantedBy=multi-user.target
